# watch_puppet.conf
[program:watch_puppet]
command=/watch_puppet.sh
